# .old/linux.Mar2022/conf.d/xterm.conf
[program:xterm]
command=rxvt
autorestart=true
